# deploy/azure-vm/systemd/python-sql-repair-model.service
[Unit]
Description=MLX SQL Repair Model Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/python-sql-parser
Environment=MODEL_ID=mlx-community/Qwen2.5-Coder-1.5B-Instruct-4bit
ExecStart=/opt/homebrew/bin/python3.11 -m mlx_lm.server --model ${MODEL_ID} --host 127.0.0.1 --port 8001 --max-tokens 1536 --log-level INFO
Restart=always
RestartSec=5
User=azureuser
Group=azureuser

[Install]
WantedBy=multi-user.target

